# ispider.conf
[program:ineww]
command=scrapyd
directory=/app/
stdout_logfile=/var/log/ispider/supervisor_stdout.log
stderr_logfile=/var/log/ispider/supervisor_stderr.log
autostart=true
autorestart=true
